# openpilot_installer/rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Apache2..."
exec httpd -D FOREGROUND