# rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting PHP-FPM..."
exec php-fpm8 -F